# config/systemd/user/meshtastic-notifyd.service
[Unit]
Description=Meshtastic notification daemon
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/meshtastic-notifyd
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
