# udev/hdmi.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/sbin/runuser -u krim -- /home/krim/.config/sway/hdmisound.sh"
